# home/.config/autostart/com.bitwarden.desktop.desktop
[Desktop Entry]
Type=Application
Name=com.bitwarden.desktop
X-XDP-Autostart=com.bitwarden.desktop
Exec=flatpak run --command=bitwarden.sh com.bitwarden.desktop --autostart %U
X-Flatpak=com.bitwarden.desktop
X-GNOME-Autostart-enabled=false
